# backup/backup.crontab
0 3 * * * bash /home/hector/server/backup/backup.sh >> /home/hector/server/backup/backup.log 2>&1
